# scripts/systemd/computehub-backup.timer
[Unit]
Description=ComputeHub backup harian (workspace + konfigurasi)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
